# luck-agent.service
[Unit]
Description=Luck Agent V2 (Goal Runtime)
After=network.target

[Service]
Type=simple
User=luck-agent
Group=luck-agent
WorkingDirectory=/opt/luck-agent
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/opt/luck-agent/.env
ExecStart=/opt/luck-agent/venv/bin/python main.py
Restart=always
RestartSec=5
NoNewPrivileges=no
ProtectSystem=strict
StandardOutput=append:/var/log/luck-agent.log
StandardError=append:/var/log/luck-agent.log

[Install]
WantedBy=multi-user.target
